# scshafe-test-service.service

[Unit]
Description=scshafe test of systemd services

[Service]
ExecStart=/usr/bin/scshafe-test-service
Type=simple
